# hermes-lean-backup/systemd-user/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway - Messaging Platform Integration
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/rapiduser/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run --replace
WorkingDirectory=/rapiduser/.hermes/hermes-agent
Environment="PATH=/rapiduser/.hermes/hermes-agent/venv/bin:/rapiduser/.hermes/hermes-agent/node_modules/.bin:/rapiduser/.hermes/node/bin:/home/rapiduser/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/rapiduser/.hermes/hermes-agent/venv"
Environment="HERMES_HOME=/rapiduser/.hermes"
Restart=always
RestartSec=5
RestartMaxDelaySec=300
RestartSteps=5
RestartForceExitStatus=75
KillMode=mixed
KillSignal=SIGTERM
ExecReload=/bin/kill -USR1 $MAINPID
TimeoutStopSec=210
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
